# puppet/module/templates/initd-script.erb
#!/bin/bash
# <%= @service_name %>        Init script for <%= @service_name %>
# chkconfig: 2345 98 02
# description: <%= @service_name %>
# processname: <%= @service_name %>

PIDFILE="<%= @pid_file %>"

LOG_DIR="<%= @log_dir %>"
startscript="<%= @start_file %>"

if [ -r "/etc/sysconfig/<%= @service_name -%>" ]; then
    . /etc/sysconfig/<%= @service_name %>
fi

export LOG_DIR

# Source function library.
. /etc/init.d/functions

start() {
    # Start daemons.
    echo -n $"Starting <%= @service_name %> "
    daemon --user=<%= @user %> --pidfile=$PIDFILE $startscript

    RETVAL=$?

    pgrep -f -u <%= @user %> <%= @binary_file %> > $PIDFILE
    chown <%= @user %>:<%= @group %> $PIDFILE

    if [[ $RETVAL != 0 ]]; then
        failure
        exit $RETVAL
    fi

    echo
}

stop() {
	if [ ! -f $PIDFILE ]; then
		echo -n "<%= @service_name %> not running. "
		success
		echo
		return 0
	fi
    echo -n $"Shutting down <%= @service_name %>: "
	killproc -p $PIDFILE
	[[ -f $PIDFILE ]] && rm -f $PIDFILE
	RETVAL=$?
        echo
	return $RETVAL
}

do_status() {
    status -p $PIDFILE "<%= @service_name %>"
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  status)
    do_status
	;;
  restart|reload|force-reload)
	stop
	start
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|force-reload}"
	exit 2
esac
